# dist/redhat/systemd/scylla-jmx.service
[Unit]
Description=Scylla JMX
Requires=scylla-server.service
After=scylla-server.service

[Service]
Type=simple
EnvironmentFile=/etc/sysconfig/scylla-jmx
User=scylla
Group=scylla
ExecStart=/usr/lib/scylla/jmx/scylla-jmx -l /usr/lib/scylla/jmx
KillMode=process
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
