# dist/ubuntu/debian/scylla-jmx.upstart
# scylla-jmx - ScyllaDB
#
# ScyllaDB

description	"ScyllaDB jmx"

start on started scylla-server
stop on stopping scylla-server

umask 022

console log

script
    . /etc/default/scylla-jmx
    export JMX_LOCAL_PORT
    exec /usr/lib/scylla/jmx/jmx_run
end script
